# edison-image-edison-ext4/var/lib/opkg/info/alsa-utils-alsactl.list
/lib/systemd/system/alsa-state.service
/usr/share/alsa/init/ca0106
/usr/share/alsa/init/test
/lib/udev/rules.d/90-alsa-restore.rules
/lib/systemd/system/basic.target.wants/alsa-restore.service
/usr/share/alsa/init/info
/lib/systemd/system/alsa-restore.service
/usr/share/alsa/init/hda
/lib/systemd/system/shutdown.target.wants/alsa-store.service
/usr/share/alsa/init/default
/var/lib/alsa/
/lib/systemd/system/shutdown.target.wants/
/lib/systemd/system/basic.target.wants/
/usr/share/alsa/init/00main
/usr/share/alsa/init/
/usr/share/alsa/init/help
/lib/systemd/system/alsa-store.service
/lib/systemd/system/basic.target.wants/alsa-state.service
/usr/sbin/alsactl
/var/lib/alsa/asound.state
