# edison-image-edison-ext4/var/lib/opkg/info/ap-mode-toggle.list
/lib/systemd/system/ap-mode-toggle.service
/usr/bin/ap-mode-toggle
