# edison-image-edison-ext4/var/lib/opkg/info/base-files.list
/mnt/
/usr/src/
/etc/systemd/system/local-fs.target.wants/
/etc/nsswitch.conf
/var/lock
/etc/host.conf
/var/tmp
/etc/shells
/etc/motd
/var/run
/etc/issue
/etc/systemd/system/default.target.wants/factory.mount
/etc/systemd/system/local-fs.target.wants/media-sdcard.automount
/sys/
/var/local/
/etc/fstab
/etc/inputrc
/lib/systemd/system/media-sdcard.mount
/etc/systemd/system/default.target.wants/
/usr/share/man/
/media/
/etc/systemd/
/usr/games/
/etc/release
/home/
/etc/skel/.bashrc
/lib/systemd/system/factory.mount
/usr/share/misc/
/bin/
/tmp/
/etc/skel/
/etc/hostname
/var/spool/
/etc/systemd/system/
/usr/sbin/
/var/lib/misc/
/var/volatile/tmp/
/dev/
/var/volatile/
/var/volatile/log/
/proc/
/run/
/etc/skel/.profile
/var/backups/
/etc/profile
/usr/share/dict/
/etc/default/usbd
/etc/filesystems
/etc/mtab
/lib/systemd/system/media-sdcard.automount
/etc/issue.net
/usr/share/info/
/home/root/
