# edison-image-edison-ext4/var/lib/opkg/info/blink-led.list
/lib/systemd/system/blink-led.service
/usr/bin/blink-led
