# edison-image-edison-ext4/var/lib/opkg/info/bluetooth-rfkill-event.list
/etc/firmware/bcm43341.conf
/lib/systemd/system/bluetooth-rfkill-event.service
/usr/sbin/bluetooth_rfkill_event
