# edison-image-edison-ext4/var/lib/opkg/info/bluetooth-rfkill-event.postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

if type systemctl >/dev/null 2>/dev/null; then
	systemctl $OPTS enable bluetooth-rfkill-event.service

	if [ -z "$D" -a "enable" = "enable" ]; then
		systemctl restart bluetooth-rfkill-event.service
	fi
fi
