# edison-image-edison-ext4/var/lib/opkg/info/cleanjournal.list
/usr/sbin/clean_journal.sh
/lib/systemd/system/cleanjournal.service
