# edison-image-edison-ext4/var/lib/opkg/info/connman-init-systemd.list
/lib/systemd/system/connman-init.service
/usr/lib/connman/connman-init-settings.sh
/usr/lib/connman/
/usr/lib/connman/connman.settings.template
