# edison-image-edison-ext4/var/lib/opkg/info/crashlog.list
/lib/systemd/system/crashlog.service
/usr/sbin/retrieve_crashlog.sh
