# edison-image-edison-ext4/var/lib/opkg/info/dbus-1.list
/usr/bin/dbus-send
/var/
/usr/bin/dbus-uuidgen
/lib/systemd/system/dbus.socket
/usr/lib/dbus/dbus-daemon-launch-helper
/lib/systemd/system/sockets.target.wants/dbus.socket
/usr/bin/dbus-daemon
/lib/systemd/system/dbus.target.wants/dbus.socket
/etc/dbus-1/session.d/
/usr/lib/dbus/
/usr/bin/dbus-cleanup-sockets
/var/lib/dbus/
/var/lib/
/usr/bin/dbus-launch
/lib/systemd/system/dbus.service
/lib/systemd/system/multi-user.target.wants/dbus.service
/etc/default/volatiles/99_dbus
/usr/share/dbus-1/services/
/etc/dbus-1/session.conf
/lib/systemd/system/dbus.target.wants/
/etc/dbus-1/system.conf
/etc/default/volatiles/
/etc/default/
/lib/systemd/system/sockets.target.wants/
/usr/bin/dbus-monitor
/lib/systemd/system/multi-user.target.wants/
