# edison-image-edison-ext4/var/lib/opkg/info/hostapd-daemon.list
/usr/sbin/hostapd
/lib/systemd/system/hostapd.service
/usr/sbin/hostapd_cli
/etc/hostapd/hostapd.conf
/lib/systemd/system/udhcpd-for-hostapd.service
/etc/hostapd/
/etc/hostapd/udhcpd-for-hostapd.conf
