# edison-image-edison-ext4/var/lib/opkg/info/mcu-fw-load.list
/lib/systemd/system/mcu_fw_loader.service
/etc/intel_mcu/
/etc/intel_mcu/mcu_fw_loader.sh
