# edison-image-edison-ext4/var/lib/opkg/info/post-install.list
/etc/systemd/system/basic.target.wants/post-install.service
/sbin/post-install.sh
/lib/systemd/system/post-install.service
