# edison-image-edison-ext4/var/lib/opkg/info/pwr-button-handler.list
/lib/systemd/system/pwr-button-handler.service
/usr/bin/pwr_button_handler
