# edison-image-edison-ext4/var/lib/opkg/info/resize-rootfs.list
/lib/systemd/system/resize-rootfs.service
