# edison-image-edison-ext4/var/lib/opkg/info/run-postinsts.list
/lib/systemd/system/run-postinsts.service
/usr/sbin/run-postinsts
