# edison-image-edison-ext4/var/lib/opkg/info/systemd-binfmt.list
/etc/binfmt.d/
/lib/systemd/system/shutdown.target
/lib/systemd/system/proc-sys-fs-binfmt_misc.mount
/lib/systemd/system/systemd-binfmt.service
/lib/systemd/systemd-binfmt
/lib/systemd/system/proc-sys-fs-binfmt_misc.automount
/usr/lib/binfmt.d/
