# edison-image-edison-ext4/var/lib/opkg/info/systemd-serialgetty.list
/etc/systemd/system/getty.target.wants/serial-getty@ttyMFD2.service
/lib/systemd/system/serial-getty@.service
