# edison-image-edison-ext4/var/lib/opkg/info/watchdog-sample.list
/usr/bin/watchdog-sample
/lib/systemd/system/watchdog-sample.service
